# pttai-next-replicator.service
[Unit]
Description=PTT.ai-next-Replicator
After=network.target

[Service]
ExecStart=node replicator.js
WorkingDirectory=/home/ubuntu/pttai-next/backend
Type=simple
User=ubuntu
Restart=always
TimeoutSec=15

[Install]
WantedBy=multi-user.target